# scheduled-tasks/packer-updates.timer
[Unit]
Description=Run safe packer updates daily

[Timer]
OnBootSec=45min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
