# scheduled-tasks/packer-updates.service
[Unit]
Description=Run safe packer updates

[Service]
Type=simple
User=kiki
ExecStart=/path/to/safe_updates.py packer
