# scheduled-tasks/pacman-updates.service
[Unit]
Description=Run safe pacman updates

[Service]
Type=simple
ExecStart=+/path/to/safe_updates.py pacman
