# scheduled-tasks/root-daily.timer
[Unit]
Description=Run root scripts daily

[Timer]
OnBootSec=17min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
