# scheduled-tasks/root-daily.service
[Unit]
Description=Run root's daily scripts

[Service]
Type=simple
ExecStart=+/usr/local/bin/root-daily.sh
