# scheduled-tasks/root-hourly.timer
[Unit]
Description=Run root scripts hourly

[Timer]
OnBootSec=23min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
