# scheduled-tasks/root-daily.service
[Unit]
Description=Run root's daily scripts

[Service]
Type=simple
ExecStart=+/home/katie/linux-config/scheduled-tasks/daily-daily.sh
