# scheduled-tasks/root-hourly.service
[Unit]
Description=Run root's hourly scripts

[Service]
Type=simple
ExecStart=+/home/katie/linux-config/scheduled-tasks/root-hourly.sh
