# scheduled-tasks/root-hourly.service
[Unit]
Description=Run root's hourly scripts

[Service]
Type=simple
ExecStart=+/usr/local/bin/root-hourly.sh
